# PACK-warehouse/tools/com.vkoffee.warehouse-full-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vkoffee.warehouse-full-loop</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/alexander/Github/VK-offee/PACK-warehouse/tools/warehouse_full_loop.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/alexander/Github/VK-offee</string>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/alexander/Github/VK-offee/.github/scripts/sync.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/alexander/Github/VK-offee/.github/scripts/sync.error.log</string>
</dict>
</plist>
